# home/Library/LaunchAgents/danemacmillan.etc.rclone.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>KeepAlive</key>
  <true/>
  <key>Label</key>
  <string>danemacmillan.etc.rclone</string>
  <key>ProgramArguments</key>
  <array>
		<string>/usr/local/bin/rclone</string>
		<string>mount</string>
		<string>gd:</string>
		<string>/Users/danemacmillan/mnt/gd</string>
		<string>--allow-other</string>
		<string>--dir-cache-time=960h</string>
		<string>--poll-interval=15s</string>
		<string>--umask=002</string>
		<string>--vfs-cache-mode=full</string>
		<string>--vfs-read-ahead=1G</string>
		<string>--vfs-cache-max-age=480h</string>
		<string>--vfs-cache-max-size=140G</string>
		<string>--cache-dir=/Users/danemacmillan/.cache/rclone</string>
		<string>--bwlimit=4M:25M</string>
		<string>--stats=1s</string>
		<string>--log-level=INFO</string>
		<string>--log-file=/Users/danemacmillan/tmp/rclone-mount-gd.log</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>

